# etc/supervisor/conf.d/eva-onboot.conf
[program:eva-onboot]
command=/opt/eva/sbin/eva-control start
autorestart=false
autostart=true
priority=500
startsecs=0
